# initScripts/FanoutPlusPlus.gentoo
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting FanoutPlusPlus"
	start-stop-daemon --start -b -m -u fanoutplusplus:fanoutplusplus --pidfile "/var/run/FanoutPlusPlus.pid" \
		--exec FanoutPlusPlus -- "${OPTIONS}"
	eend $?
}

stop() {
	ebegin "Stopping FanoutPlusPlus"
	start-stop-daemon --stop --pidfile "/var/run/FanoutPlusPlus.pid"
	eend $?
}
